# crontab
# Run every day at 02:00 (AM)
0 2 * * * bash /srv/run-import.sh

# Run every 3 months at 00:00 (AM)
0 0 1 */3 * /usr/local/bin/php /srv/bin/console app:import-geo-addresses
