# crontab
# Run every day at 02:00 (AM)
0 2 * * * bash /srv/run-import.sh
